# roles/nas/files/backup-push.cron
0       4       *       *       *       root    /usr/local/sbin/syncBackup silicon.collegiumv.org
