# roles/xbps-mini-repo/files/xbps-mini-repo.crontab
0 0 * * * mini-repo /opt/xbps-mini-repo/xbps-mini-builder 2>&1 > /dev/null
