# roles/openbsd/files/etc/rc.d/dhcpd
#!/bin/sh

daemon="/usr/local/sbin/dhcpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	touch /var/db/dhcpd.leases
}

rc_cmd $1
