# roles/etherwake/files/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command
* * * * * root /usr/local/sbin/wakeAll /etc/macs.list
