# roles/openbsd/files/etc/rc.local
# Use the local bind server and make the domain collegiumv.org instead of DHCP information
cp /etc/resolv.conf.final /etc/resolv.conf
